# infra/systemd/signage-api.service
[Unit]
Description=OmniHub Signage API
After=network.target postgresql.service redis-server.service

[Service]
Type=simple
User=chiho
Group=chiho
WorkingDirectory=/opt/signage/apps/api
Environment=NODE_ENV=production
EnvironmentFile=/etc/signage/api.env
ExecStart=/home/chiho/.nvm/versions/node/v24.13.1/bin/node /opt/signage/apps/api/dist/index.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
